# roles/cis_security/templates/aidecheck.service
[Unit]
Description=Aide Check

[Service]
Type=simple
{% if ansible_distribution == "RedHat" or ansible_distribution == "CentOS" or ansible_distribution == "Oracle"  %}
ExecStart=/usr/sbin/aide --check
{% elif  ansible_distribution == "Ubuntu" and ansible_distribution_major_version|int >= 22  %}
ExecStart=/usr/bin/aide --check
{% else %}
ExecStart=/usr/bin/aide.wrapper -C
{% endif %}

[Install]
WantedBy=multi-user.target
