# roles/cis_security/templates/aidecheck.service
[Unit]
Description=Aide Check

[Service]
Type=simple
{% if ansible_distribution == "RedHat" or ansible_distribution == "CentOS" %}
ExecStart=/usr/sbin/aide --check
{% else %}
ExecStart=/usr/bin/aide.wrapper -C
{% endif %}

[Install]
WantedBy=multi-user.target
